# themes/commotion/root/etc/uci-defaults/luci-theme-commotion
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Commotion=/luci-static/commotion
	set luci.main.mediaurlbase=/luci-static/commotion
	commit luci
EOF
exit 0
